# src/s6/etc/s6-overlay/s6-rc.d/websocket/run
#!/command/execlineb -P
with-contenv
wait-for-it -t 120 ${REDIS_QUEUE:-redis-queue}:6379
s6-notifyoncheck
/usr/bin/node /home/frappe/frappe-bench/apps/frappe/socketio.js
