# src/s6/etc/s6-overlay/s6-rc.d/websocket/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
/usr/bin/node /home/frappe/frappe-bench/apps/frappe/socketio.js
